# ansible/roles/sunshine/templates/sunshine.service.j2
[Unit]
Description=Sunshine Game Streaming Server
Documentation=https://docs.lizardbyte.dev/projects/sunshine/
BindsTo=graphical-session.target
Wants=sway.service
After=sway.service

[Service]
Type=simple
ExecStart=/usr/bin/sunshine {{ sunshine_config_dir }}/sunshine.conf
Restart={{ 'on-failure' if sunshine_restart_on_failure else 'no' }}
RestartSec={{ sunshine_restart_sec }}

# Environment for Wayland
Environment=XDG_RUNTIME_DIR=/run/user/{{ streamin_uid }}
Environment=WAYLAND_DISPLAY={{ sunshine_display }}
Environment=XDG_CONFIG_HOME={{ streamin_home }}/.config
Environment=XDG_DATA_HOME={{ streamin_home }}/.local/share

# AMD GPU environment
Environment=LIBVA_DRIVER_NAME=radeonsi
Environment=VDPAU_DRIVER=radeonsi

[Install]
WantedBy=graphical-session.target
